# build/supervisor/live-monitor.ini
[program:live-monitor]
command=/usr/local/bin/supervisor /usr/home/jvazquez/live-monitor/app.js
directory=/usr/home/jvazquez/live-monitor/
user=jvazquez
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/live-monitor.stdout.log
stdout_logfile_maxbytes=2MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=5     ; # of stdout logfile backups (default 10)
stderr_logfile=/var/log/live-monitor.stderr.log
stderr_logfile_maxbytes=2MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=5     ; # of stderr logfile backups (default 10)

