# build/supervisor/live-monitor.conf
[program:live-monitor]
#I'm guessing you actually read this and installed supervisor sudo npm install -g supervisor ?
command=/usr/local/bin/supervisor /path/to/app.js
directory=/path/to/directory/that/contains/app.js/with/trailing/slash/
user=an-user-to-run-daemon
autostart=true
autorestart=true
redirect_stderr=true
